# net-analyzer/cyanite/files/initd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

export JAVA_HOME=`java-config --jre-home`

checkconfig() {
  if [ ! -r "${CYANITE_CONFIG_FILE}" ]; then
    eerror "CYANITE_CONFIG_FILE doesn't point to a readable config file"
  fi
}

depend() {
  need net
}

start() {
  checkconfig || return 1

  [ -e `dirname "$PIDFILE"` ] || \
    install -d -o${CYANITE_USER} -g${CYANITE_GROUP} -m750 `dirname $PIDFILE`

  start-stop-daemon \
    --start \
    --exec "${JAVA_HOME}/bin/java" \
    --user ${CYANITE_USER} \
    --group ${CYANITE_GROUP} \
    --pidfile "${CYANITE_PID_FILE}" \
    -- -jar /usr/share/cyanite/lib/cyanite.jar -f "${CYANITE_CONFIG_FILE}" ${CYANITE_OPTIONS}
}

stop() {
  start-stop-daemon \
    --stop \
    --pidfile "${CYANITE_PID_FILE}"
}
